# supervisord_flask.conf
[program:flask]
directory = /usr/src/app
command = python3 main.py
